# supervisor/dashboard.conf
[program:realtime-dashboard]
directory=/var/www/realtime-sensor-monitoring
command=/usr/bin/npm start -- --hostname 0.0.0.0 --port 3000
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/realtime-dashboard.err.log
stdout_logfile=/var/log/realtime-dashboard.out.log
environment=NODE_ENV="production",PORT="3000"